# deploy/butchr-health.service
# butchr health watchdog — one-shot probe driven by butchr-health.timer.
#
# Curls butchr's /health and restarts butchr.service if it is unreachable or
# reports a degraded/stale dispatcher tick. See scripts/health-watchdog.sh.
#
# TEMPLATE: @REPO_DIR@ is substituted by scripts/install-service.sh.
[Unit]
Description=butchr health watchdog (restart butchr if /health is unreachable or stale)
After=butchr.service

[Service]
Type=oneshot
Environment=PATH=%h/.bun/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
# Override the probe target / unit here if you run on a non-default port:
#   Environment=BUTCHR_HEALTH_URL=http://127.0.0.1:47800/health
#   Environment=BUTCHR_SERVICE=butchr.service
ExecStart=@REPO_DIR@/scripts/health-watchdog.sh
SyslogIdentifier=butchr-health
